# ft_services/ftps/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ftps]

